# wp-linksurfer.timer
[Unit]
Description=Timer for LinkSurfer Wordpress Plugin

[Timer]
OnBootSec=300
OnUnitActiveSec=1w

[Install]
WantedBy=multi-user.target